# audit_sudo_wrappers/FreeBSD/bin/crontab
#!/bin/sh

if [ -f /usr/bin/crontab ];
	then /usr/local/bin/sudo /usr/bin/crontab "$@"
fi
