# 8/audit_sudo_wrappers/audit_sudo_wrappers/RED OS/bin/crontab
#!/usr/bin/sh

if [ -f /usr/bin/crontab ];
        then /usr/bin/sudo /usr/bin/crontab "$@"
fi
